# gateway/basestation.service
[Unit]
Description=IoT Security Base Station Gateway
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
WorkingDirectory=/home/pi/iot-security/gateway
ExecStart=/home/pi/iot-security/gateway/.venv/bin/python gateway_service.py
Restart=always
RestartSec=10
User=pi
Group=pi

# Load environment variables
EnvironmentFile=/home/pi/iot-security/gateway/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=basestation

[Install]
WantedBy=multi-user.target
